# automation-system/com.brandedflow.morning.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.brandedflow.morning</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/jencortez-walters/brandedflow/automation-system/morning_routine_v2.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>StandardOutPath</key>
    <string>/Users/jencortez-walters/brandedflow/automation-system/logs/morning.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jencortez-walters/brandedflow/automation-system/logs/morning.error.log</string>
</dict>
</plist>
